# etc/systemd/system/android-mount.service.d/10-schedtune.conf
[Service]
ExecStartPre=-/usr/bin/umount -l /sys/fs/cgroup/schedtune
