# etc/systemd/system/android-mount.service.d/20-mtkdispmgr.service
[Service]
ExecStart=-/usr/bin/chown droidian:droidian /dev/mtk_disp_mgr
