# etc/systemd/user/pulseaudio.service.d/10-env.conf
[Service]
Environment="HYBRIS_LD_LIBRARY_PATH=/apex/com.android.runtime/lib64:/system/lib64:/vendor/lib64"
