# usr/lib/systemd/user/pulseaudio.service.d/10-env.service
[Service]
Environment="HYBRIS_LD_LIBRARY_PATH=/apex/com.android.runtime/lib:/system/lib:/vendor/lib"
